# systemd/lightning-eye.service
[Unit]
Description=Lightning Eye Blitzsensor
After=network-online.target graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
WorkingDirectory=%h/lightning_eye
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
ExecStart=%h/lightning_eye/start.sh
ExecStopPost=%h/lightning_eye/.venv/bin/python -m app.leds_off
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
